# overlay/systemd/send-status-update.service
[Unit]
Description=Send device status update

[Service]
Type=oneshot
ExecStart=/sensos/libexec/send-status-update.sh
