# overlay/systemd/sensos-birdnet.service
[Unit]
Description=Run BirdNET over queued WAV files and emit FLAC label runs
After=multi-user.target sensos-arecord.service
StartLimitIntervalSec=300
StartLimitBurst=5

# This service is optional. It also requires:
# - BirdNET model files under /sensos/birdnet/BirdNET_v2.4_tflite
# - host-native Python deps from python/requirements-birdnet.txt

[Service]
Type=simple
User=sensos-runner
ExecStart=/sensos/python/birdnet-venv/bin/python /sensos/libexec/process-birdnet.py
WorkingDirectory=/sensos/data
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
